# system_files/shared/etc/profile.d/99-greatlem0n-fastfetch.sh
# shellcheck shell=bash
unalias neofetch 2>/dev/null
unalias neowofetch 2>/dev/null
unalias fastfetch 2>/dev/null

alias neofetch='/usr/libexec/greatlem0n-fastfetch'
alias neowofetch='/usr/libexec/greatlem0n-fastfetch'
alias fastfetch='/usr/libexec/greatlem0n-fastfetch'
